# upstart/node-ownpt.conf
start on (started solr)
stop on runlevel [016]
setuid ibmcloud
kill timeout 30
respawn
env API_KEY="<INSERT API KEY HERE>"
script
  chdir /home/ibmcloud/ownpt2
  exec nodejs app.js
end script
