# upstart/solr.conf
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
setuid solr
kill timeout 30
script
  chdir /usr/local/solr
  exec bin/solr start -f
end script
pre-stop script
  chdir /usr/local/solr
  exec bin/solr stop
end script
